# .config/systemd/user/active_window.timer
[Unit]
Description=Log Battery Percentage for debugging

[Timer]
Persistent=true
OnCalendar=*:*:0/30

[Install]
WantedBy=timers.target


