# .config/systemd/user/aptmanual.timer
[Unit]
Description=Save apt-mark showmanual

[Timer]
Persistent=true
#OnUnitActiveSec=1day
OnCalendar=00:00:00


[Install]
WantedBy=timers.target
