# .config/systemd/user/bitwarden.service
[Unit]
Description=Password Manager
After=default.target

[Service]
ExecStart=/usr/bin/env bitwarden

[Install]
WantedBy=default.target
