# .config/systemd/user/documents-backup.service
[Unit]
Description=Documents Backup

[Service]
Type=simple
ExecStart=/bin/sh -c "rclone sync ~/Documents mega:Documents"


