# .config/systemd/user/dunst.service
[Unit]
Description=Notification daemon
After=default.target

[Service]
ExecStart=/usr/bin/env dunst

[Install]
WantedBy=default.target
