# .config/systemd/user/dunst.service
[Unit]
Description=Notification daemon
PartOf=graphical.target

[Service]
ExecStart=/usr/bin/env dunst

[Install]
WantedBy=graphical.target
