# .config/systemd/user/gitbu.service
[Unit]
Description=Git Backup for Ignored files

[Service]
Type=simple
ExecStart=%h/.local/lib/scripts/gitbu.sh

