# .config/systemd/user/graphical.target
[Unit]
Description=launchable graphical-session.target
BindsTo=graphical-session.target
