# .config/systemd/user/libinput-gestures.service
[Unit]
Description=Touchpad Gestures
After=default.target

[Service]
ExecStart=/usr/bin/env libinput-gestures

[Install]
WantedBy=default.target
