# .config/systemd/user/logpercentage.service
[Unit]
Description=Log Battery Percentage for debugging

[Service]
ExecStart=%h/.config/systemd/user/logpercentage.sh

