# .config/systemd/user/rclone-sync@.service
[Unit]
Description=Rclone sync for %i

[Service]
Type=simple
ExecStart=/usr/bin/zsh -c "source %E/rclone/sync.sh %i; start"
