# .config/systemd/user/redshift.service
[Unit]
Description=Redshift Screen Dimmer
After=default.target

[Service]
ExecStart=/usr/bin/env redshift-gtk

[Install]
WantedBy=default.target
