# .config/systemd/user/x11vnc.service
[Unit]
Description=Start x11vnc at startup.
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/x11vnc -auth guess -forever -noxdamage -repeat -rfbauth .vnc/passwd -rfbport 5900 -shared

[Install]
WantedBy=multi-user.target
